# ansible/roles/bluetooth/templates/bt-speaker.service.j2
[Unit]
Description=Bluetooth Speaker Auto Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bt-speaker-agent.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
